# FILES/cron/script-3.service
[Unit]
Description=Dropbear User-Limit Service

[Service]
Type=simple
ExecStart=/usr/bin/script-3

[Install]
WantedBy=multi-user.target